# ansible/roles/self_pull/templates/foundry-pull.service.j2
# {{ ansible_managed }}
[Unit]
Description=Foundry self-managed ansible-pull converge
Documentation=file://{{ self_pull_checkout_dir }}/ansible/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User={{ self_pull_run_user }}
Group={{ self_pull_run_group }}
ExecStart={{ self_pull_wrapper_path }}
# Journald captures stdout/stderr; keep it the default.
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=5
