# roles/ansible-pull/templates/ansible-pull.service.j2
[Unit]
Description=Run foundry ansible-pull
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart={{ ansible_pull_wrapper_path }}
